# templates/etc/network/if-pre-up.d/configure_br-ex
#!/bin/sh
ip addr del {{ external_network }} dev eth0
ip link set eth0 up
ip link set eth0 promisc on
ip addr add {{ external_network }} dev br-ex
ip link set br-ex up