# templates/etc/network/if-pre-up.d/configure_br-ex
#!/bin/sh
ip addr del {{ CONTROLLER_EXT_IP }}/16 dev eth3
ip link set eth3 up
ip link set eth3 promisc on
ip addr add {{ CONTROLLER_EXT_IP }}/16 dev br-ex
ip link set br-ex up